# ide/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: IDE
# Runs the nginx daemon
# ==============================================================================

# Wait for c9 to become available
bashio::net.wait_for 80

bashio::log.info 'Starting Nginx...'

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

if bashio::config.true 'ssl'; then
    exec nginx -c /etc/nginx/nginx-ssl.conf
else
    exec nginx -c /etc/nginx/nginx.conf
fi
