# ide/rootfs/etc/services.d/c9/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: IDE
# Runs the cloud9 server
# ==============================================================================
declare -a c9_options

bashio::log.info 'Starting the C9...'

# Set daemon debug level
if bashio::debug; then
    c9_options+=(--debug)
fi

# Add default options
c9_options+=(--port 80)
c9_options+=(-w '/workspace')
c9_options+=(--workspacetype 'homeassistant')

# Change working directory
cd /cloud9 || bashio::exit.nok 'Unable to change working directory'

# Run the c9 server
exec node /cloud9/server.js "${c9_options[@]}"
